# bot_manager.service
[Unit]
Description=Bot Manager

[Service]
Type=simple
ExecStart=/root/BotManager/BotManager

[Install]
WantedBy=multi-user.target