# web_bot.service
[Unit]
Description=Bot Manager

[Service]
Type=simple
ExecStart=/root/web_bot/WebAccessibleBot

[Install]
WantedBy=multi-user.target